# debian/dibbler-relay.postrm
#!/bin/sh
set -e
# Automatically added by dh_installinit
if [ "$1" = "purge" ] ; then
	update-rc.d dibbler-relay remove >/dev/null || exit 0
fi
# End automatically added section
