# debian/dibbler-server.prerm
#!/bin/sh
set -e

#if [ -x "/etc/init.d/dibbler-server" ]; then
#	if [ -x "`which invoke-rc.d 2>/dev/null`" ]; then
#		invoke-rc.d dibbler-server stop || exit 0
#	else
#		/etc/init.d/dibbler-server stop || exit 0
#	fi
#fi

#DEBHELPER#
