# debian/dibbler-client.prerm
#!/bin/sh
set -e

#if [ -x "/etc/init.d/dibbler-client" ]; then
#	if [ -x "`which invoke-rc.d 2>/dev/null`" ]; then
#		invoke-rc.d dibbler-client stop || exit 0
#	else
#		/etc/init.d/dibbler-client stop || exit 0
#	fi
#fi

#DEBHELPER#

